# cloud/systemd/portfolio-ai-tunnel.service
[Unit]
Description=SSH Tunnel to T5810 (vLLM + Qdrant + embeddings + reranker for Portfolio AI)
After=network.target sshd.service
Wants=network-online.target

[Service]
Type=simple
User=root
# Forwards localhost ports on the VPS to the T5810's LAN-only services so the
# api-proxy can reach them at 127.0.0.1. The home server's inference services
# are never exposed to the internet — only reachable through this tunnel.
#   8001=reserved  8004=vLLM  8005=embeddings  8006=reranker  6333=Qdrant  8788=headroom-lib
#
# Per-machine config lives in /etc/default/portfolio-ai-tunnel (NOT in the repo):
#   T5810_USER=chris
#   T5810_HOST=ai.cwetzel.com    # DNS preferred over a raw IP; survives Verizon dynamic-IP rotation
# Create the file on the VPS once with `T5810_USER` and `T5810_HOST` set,
# then `systemctl daemon-reload && systemctl restart portfolio-ai-tunnel.service`.
EnvironmentFile=/etc/default/portfolio-ai-tunnel
ExecStart=/usr/bin/ssh -N -o StrictHostKeyChecking=accept-new -o ServerAliveInterval=60 -o ServerAliveCountMax=3 -L 127.0.0.1:8001:127.0.0.1:8001 -L 127.0.0.1:8004:127.0.0.1:8004 -L 127.0.0.1:8005:127.0.0.1:8005 -L 127.0.0.1:8006:127.0.0.1:8006 -L 127.0.0.1:6333:127.0.0.1:6333 -L 127.0.0.1:8788:127.0.0.1:8788 ${T5810_USER}@${T5810_HOST}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
