# cloud/systemd/portfolio-ai-tunnel.service
[Unit]
Description=SSH Tunnel to T5810 (vLLM + Qdrant + embeddings + reranker for Portfolio AI)
After=network.target sshd.service
Wants=network-online.target

[Service]
Type=simple
User=root
# Forwards localhost ports on the VPS to the T5810's LAN-only services so the
# api-proxy can reach them at 127.0.0.1. The home server's inference services
# are never exposed to the internet — only reachable through this tunnel.
#   8001=reserved  8004=vLLM  8005=embeddings  8006=reranker  6333=Qdrant  8788=headroom-lib
#   8007=faithfulness verifier — on the asrock B550 (home LAN), NOT the T5810.
#
# Verifier (verifier-faithfulness-layer.md §7.2): the verifier runs on the asrock B550
# (asrock:8007). The T5810 terminates this tunnel and shares the LAN with asrock, so
# the -L 8007 forward below targets asrock's IP and the T5810 routes to it (confirmed
# T5810 -> asrock:8007 reachable over the LAN). Then VERIFIER_URL=http://127.0.0.1:8007 on the
# api-proxy unit makes the proxy's fire-and-forget /verify reach asrock.
#
# Per-machine config lives in /etc/default/portfolio-ai-tunnel (NOT in the repo):
#   T5810_USER=chris
#   T5810_HOST=ai.cwetzel.com    # DNS preferred over a raw IP; survives Verizon dynamic-IP rotation
#   VERIFIER_HOST=<asrock-lan-ip>  # asrock's LAN address, as resolved FROM the T5810
# Create the file on the VPS once with all three set, then
# `systemctl daemon-reload && systemctl restart portfolio-ai-tunnel.service`.
#
# NOTE: VERIFIER_HOST was previously a literal LAN address in this file. It is parameterized so
# the public repo carries no internal addressing. The unit currently installed on the VPS still
# has the literal; adopting this version requires adding VERIFIER_HOST to /etc/default FIRST —
# an empty value yields `-L 127.0.0.1:8007::8007`, and ssh will refuse to start the tunnel.
EnvironmentFile=/etc/default/portfolio-ai-tunnel
ExecStart=/usr/bin/ssh -N -o StrictHostKeyChecking=accept-new -o ServerAliveInterval=60 -o ServerAliveCountMax=3 -L 127.0.0.1:8001:127.0.0.1:8001 -L 127.0.0.1:8004:127.0.0.1:8004 -L 127.0.0.1:8005:127.0.0.1:8005 -L 127.0.0.1:8006:127.0.0.1:8006 -L 127.0.0.1:6333:127.0.0.1:6333 -L 127.0.0.1:8788:127.0.0.1:8788 -L 127.0.0.1:8007:${VERIFIER_HOST}:8007 ${T5810_USER}@${T5810_HOST}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
